# crontab
*/5 * * * * /bin/unrar_torrent.sh "/data/"
